# recorder-videocap.service
[Unit]
Description=Videocap Recorder %i

[Service]
ExecStart=/opt/videocap/run.sh config-%i.sh recorder

[Install]
WantedBy=multi-user.target
